# blueberry.service
[Unit]
Description=Run the a2dp script

[Service]
ExecStart=/usr/bin/blueberry

[Install]
WantedBy=multi-user.target