# blueberry.service
[Unit]
Description=Run the a2dp auto-connecting script

[Service]
ExecStart=/usr/bin/blueberry on

[Install]
WantedBy=multi-user.target